# systemd/jetbackup-remote.service
[Unit]
Description=JetBackup Remote Orchestrator Daemon
Documentation=https://github.com/AichaDigital/jetbackup-remote
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m jetbackup_remote -c /etc/jetbackup-remote/config.json daemon
Restart=on-failure
RestartSec=300
Environment=PYTHONPATH=/opt/jetbackup-remote/src
EnvironmentFile=/etc/jetbackup-remote/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jetbackup-remote

# Security hardening (ProtectHome disabled: needs /root/.ssh for SSH keys)
ProtectSystem=strict
ReadWritePaths=/var/log /tmp /etc/jetbackup-remote /var/lib/jetbackup-remote
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target
