# systemd/jetbackup-remote.timer
[Unit]
Description=JetBackup Remote Orchestrator Timer
Documentation=https://github.com/AichaDigital/jetbackup-remote

[Timer]
# Run daily at 02:00 local time
OnCalendar=*-*-* 02:00:00
# Add random delay up to 5 minutes to avoid exact-second collisions
RandomizedDelaySec=300
# Do NOT catch up missed runs - prevents cascading concurrent backups
# if previous run overran or machine was down at scheduled time.
# The lock file (fcntl.flock) provides additional protection.
Persistent=false

[Install]
WantedBy=timers.target
